# rootfs/etc/services.d/rtl_433_mqtt_hass/run
#!/usr/bin/with-contenv bashio

declare host
declare password
declare port
declare username
declare config_path

config_path=/data/options.json

if bashio::services.available "mqtt"; then
    host=$(bashio::services "mqtt" "host")
    password=$(bashio::services "mqtt" "password")
    port=$(bashio::services "mqtt" "port")
    username=$(bashio::services "mqtt" "username")
else
    bashio::log.info "mqtt is not available"
    exit 1
fi

/usr/bin/with-contenv /usr/local/bin/python -u /usr/local/bin/rtl_433_mqtt_hass.py -u $username -P $password -H $host -p $port -r -f 2>&1
